# skills/general/autoresearch-toolfinder/systemd/autoresearch-index.timer
[Unit]
Description=Weekly refresh of autoresearch-toolfinder index
Documentation=https://github.com/jajupmochi/agent-harness

[Timer]
OnCalendar=Sun *-*-* 04:00:00
Persistent=true
RandomizedDelaySec=1h

[Install]
WantedBy=timers.target
